# images/log/supervisord.conf
[supervisord]
nodaemon=true

[program:grafana]
command=grafana-server --config=/etc/grafana/grafana.ini
autostart=true
autorestart=true
stderr_logfile=/var/log/grafana/grafana.err.log
stdout_logfile=/var/log/grafana/grafana.out.log

[program:loki]
command=loki -config.file=/etc/loki/loki-config.yaml
autostart=true
autorestart=true
stderr_logfile=/var/log/loki/loki.err.log
stdout_logfile=/var/log/loki/loki.out.log

[program:promtail]
command=promtail -config.file=/etc/promtail/promtail-config.yaml
autostart=true
autorestart=true
stderr_logfile=/var/log/promtail/promtail.err.log
stdout_logfile=/var/log/promtail/promtail.out.log
